# submissions/track2-coreline-systems-limited-openalpha-sentinel/deploy/systemd/openalpha-worker.service
[Unit]
Description=OpenAlpha Sentinel continuous collector
After=network-online.target openalpha-api.service
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/openalpha-sentinel.env
ExecStart=/usr/local/bin/openalpha-service worker
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
