# submissions/track2-coreline-systems-limited-openalpha-sentinel/deploy/systemd/openalpha-llama.service
[Unit]
Description=OpenAlpha Sentinel llama.cpp ROCm model server
After=network.target

[Service]
Type=simple
EnvironmentFile=/etc/openalpha-sentinel.env
ExecStart=/usr/local/bin/openalpha-service llama
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
