# submissions/track2-coreline-systems-limited-openalpha-sentinel/deploy/systemd/openalpha-api.service
[Unit]
Description=OpenAlpha Sentinel local API and web interface
After=network-online.target openalpha-llama.service

[Service]
Type=simple
EnvironmentFile=/etc/openalpha-sentinel.env
ExecStart=/usr/local/bin/openalpha-service api
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
